# tts-supervisor.conf
[program:tts-model]
command=/opt/tts/start-model.sh
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/var/log/portal/vllm.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/portal/vllm.log
stderr_logfile_maxbytes=0
priority=10
